# znc.conf
[program:znc]
command=su -c "/usr/bin/znc --foreground" znc
autorestart=true

