# znc.conf
[program:znc]
command=/usr/bin/znc --foreground --datadir /srv/znc
user=znc
env=HOME=/home/znc
autorestart=true
